# LaunchAgent/com.macflow.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Macflow LaunchAgent template.
  install.sh replaces __BINARY_PATH__, __LOG_OUT__ and __LOG_ERR__ with real paths
  and copies this file to ~/Library/LaunchAgents/com.macflow.agent.plist.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macflow.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
    </array>

    <!-- Starts at login and keeps it alive (restarts if it crashes). -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Interactive app: runs in the user's graphical session (Aqua). -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__LOG_OUT__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_ERR__</string>
</dict>
</plist>
